# Linux/aiknock.service
[Unit]
Description=AIKNOCK Decision Plane
Documentation=https://github.com/luigidelgiudice/AiKnock-Protocol-OS-level
After=network.target
DefaultDependencies=no
Before=basic.target

[Service]
Type=simple
ExecStart=/usr/local/bin/aiknock-svc
Restart=on-failure
RestartSec=5
RuntimeDirectory=aiknock
RuntimeDirectoryMode=0755
User=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aiknock

[Install]
WantedBy=multi-user.target
